# cookbooks/solo_client/templates/default/deregister.conf.erb

start on shutdown
task

exec chef-solo -j <%=node.delete_me_attribs%>

