# cookbooks/solo_client/templates/default/deregister.conf.erb

start on shutdown
task

exec /usr/local/bin/delete_sdb_record.py chef <%=@node.fqdn%>

