# cookbooks/solo_client/templates/default/deregister.conf.erb

start on runlevel RUNLEVEL=[016] PREVLEVEL=[2345]
task

pre-start script
    stop chef || true
end script

exec python /usr/local/bin/delete_sdb_record.py chef <%=@node.fqdn%> <%=@node.deployment%>

