# cookbooks/solo_client/templates/default/deregister.conf.erb

start on runlevel RUNLEVEL=[016] PREVLEVEL=[2345] and starting rc RUNLEVEL=[016]
task

pre-start script
    stop chef || true
end script

exec /usr/local/bin/delete_sdb_record.py chef <%=@node.fqdn%>

