# cookbooks/solo_client/templates/default/deregister.conf.erb

start on shutdown
task

pre-start script
    stop chef || true
end script

exec /usr/local/bin/delete_sdb_record.py chef <%=@node.fqdn%>

